# files/cron.d/graphite-core-temp
* * * * * pi /home/pi/rfx-commands/scripts/to-graphite-coretemp.sh
