# files/cron.d/magnets-to-graphite
*/5 * * * * pi /home/pi/rfx-commands/scripts/to-graphite-switches.sh
