# files/cron.d/nexa-lights
#
#	Lamporna
#
#00 06 * * * pi /home/pi/rfx-commands/triggers/lights.sh 00D81332 on  2
#15 08 * * * pi /home/pi/rfx-commands/triggers/lights.sh 00D81332 off 2
#
#15 17 * * * pi /home/pi/rfx-commands/triggers/lights.sh 00D81332 "group on"  1
30 23 * * * pi /home/pi/rfx-commands/triggers/lights.sh 00D81332 "group off" 1
#
#	Ebbas rum
#
#20 06 * * 1-5 pi /home/pi/rfx-commands/triggers/lights.sh 00D81332 on  3
#30 08 * * 1-5 pi /home/pi/rfx-commands/triggers/lights.sh 00D81332 off 3
#
#30 17 * * * pi /home/pi/rfx-commands/triggers/lights.sh 00D81332 on  3
00 23 * * * pi /home/pi/rfx-commands/triggers/lights.sh 00D81332 off 3
#
#	Julgranen
#
#00 06 * * * pi /home/pi/rfx-commands/triggers/lights.sh 00D81332 on  4
#15 08 * * * pi /home/pi/rfx-commands/triggers/lights.sh 00D81332 off 4
#
#25 16 * * * pi /home/pi/rfx-commands/triggers/lights.sh 00D81332 on  4
#45 23 * * * pi /home/pi/rfx-commands/triggers/lights.sh 00D81332 off 4
