# files/cron.d/nexa-lights
#
#	Lamporna
#
#00 06 * * * pi /home/pi/rfx-commands/commands/cmd-to-nexa.sh "group on"
15 08 * * * pi /home/pi/rfx-commands/commands/cmd-to-nexa.sh "group off"
#
00 18 * * * pi /home/pi/rfx-commands/commands/cmd-to-nexa.sh "group on"
30 23 * * * pi /home/pi/rfx-commands/commands/cmd-to-nexa.sh "group off"
#
#	Ebbas rum
#
15 05 * * 1-5 pi /home/pi/rfx-commands/commands/cmd-to-nexa.sh 3 on
30 07 * * 1-5 pi /home/pi/rfx-commands/commands/cmd-to-nexa.sh 3 off
#
#30 17 * * * pi /home/pi/rfx-commands/commands/cmd-to-nexa.sh 3 on
00 23 * * * pi /home/pi/rfx-commands/commands/cmd-to-nexa.sh 3 off
#
#	Julgranen
#
#00 06 * * * pi /home/pi/rfx-commands/commands/cmd-to-nexa.sh 4 on
#15 08 * * * pi /home/pi/rfx-commands/commands/cmd-to-nexa.sh 4 off
#
#45 17 * * * pi /home/pi/rfx-commands/commands/cmd-to-nexa.sh 4 on
#45 23 * * * pi /home/pi/rfx-commands/commands/cmd-to-nexa.sh 4 off
