# files/cron.d/temperatur-nu
*/5 * * * * pi /home/pi/rfx-commands/scripts/temperatur-nu.sh
