# files/cron.d/update-openhab
*/10 * * * * pi /home/pi/rfx-commands/update-openhab.sh
