# files/cron.d/update-temperatur-nu
*/5 * * * * pi /home/pi/rfx-commands/scripts/update-temperatur-nu.sh
